# deployment/systemd/invoice-extraction.service
[Unit]
Description=Invoice Extraction SQS Worker
After=network.target

[Service]
Type=simple
User=invoice-processor
WorkingDirectory=/home/ubuntu/sledge-invoice-processor/invoice-extraction
EnvironmentFile=/home/ubuntu/sledge-invoice-processor/.env
ExecStart=uv run server_handler.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target